# containers/services/server/lfg-backend.container
[Unit]
Description=LFG backend server

[Container]
ContainerName=systemd-lfg-backend

Image=ghcr.io/lookingforgrp-rit/lfg-server:latest
Pull=never

# Publish the server port
PublishPort=127.0.0.1:8081:3000

[Service]
Restart=always
EnvironmentFile=/home/woot/env/backend.env

[Install]
# Start on server boot
WantedBy=default.target